# docker/crontab
# m h dom mon dow   command
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * cd /imap2nas && npm start >> /var/log/cron.log 2>&1
